# modules/webapp/files/uwsgi.conf
description "uWSGI"
start on runlevel [2345]
stop on runlevel [06]
respawn

env USER=www-data
env UWSGI=/usr/local/bin/uwsgi
env LOGTO=/var/log/uwsgi/emperor.log
env RUN=/var/run/uwsgi

pre-start script
mkdir $RUN && chown $USER $RUN
end script

script
exec $UWSGI --master --emperor /etc/uwsgi/vassals --uid $USER --gid $USER --logto $LOGTO
end script